# profile.d/golang.sh
# begin golang.sh

pathappend /usr/local/go/bin

# end golang.sh
