# profile.d/pythondocs.sh
# Begin pythondocs.sh 

export PYTHONDOCS=/usr/share/doc/python-3/html

# End pythondocs.sh
